# scripts/upstart/ircd.js.conf
description "ircd.js - IRC server implemented in node.js"

start on filesystem
stop on runlevel [016]

exec ircdjs
